# crontab
0 * * * * /usr/local/bin/python -u /app/main.py
